# gdm.path
[Unit]

Description=Rise GDM service automatically if a graphic card exist

[Path]
PathExists=/dev/dri/card0
